# cronjobs/crontab
0 4 * * * ~/.dotfiles/cronjobs/daily_cleanup.sh