# launchd/linux/daily_cleanup.timer
[Unit]
Description=Run daily cleanup script at 2 AM

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
Unit=daily_cleanup.service  # Explicitly linking to the service

[Install]
WantedBy=timers.target
