# launchd/linux/daily_cleanup.service
[Unit]
Description=Daily cleanup
After=network.target

[Service]
Type=simple
ExecStart=/home/grado/.dotfiles/launchd/daily_cleanup.sh
